# launchagents/com.herdr.amphetamine.monitor.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>__LABEL__</string>
    <key>ProgramArguments</key>
    <array>
        <string>__PYTHON__</string>
        <string>__PLUGIN_ROOT__/scripts/monitor.py</string>
        <string>--daemon</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>__LOG_DIR__/monitor.out.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/monitor.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <!-- Bootstrap only: tunables live in config.json (TUI-editable) and are
             re-read by the daemon every poll. The plist pins the absolute herdr
             path (a user LaunchAgent has a minimal PATH) and the resolved
             per-session config/state dirs (under herdr's plugin dirs). -->
        <key>HERDR_BIN_PATH</key>
        <string>__HERDR_BIN__</string>
        <key>HERDR_SOCKET_PATH</key>
        <string>__HERDR_SOCKET__</string>
        <key>HERDR_AMPHETAMINE_CONFIG_DIR</key>
        <string>__CONFIG_DIR__</string>
        <key>HERDR_AMPHETAMINE_STATE_DIR</key>
        <string>__STATE_DIR__</string>
        <key>HERDR_AMPHETAMINE_PLIST_PATH</key>
        <string>__PLIST__</string>
        <key>HERDR_AMPHETAMINE_LABEL</key>
        <string>__LABEL__</string>
        <key>HERDR_AMPHETAMINE_AUTO_UNLOAD</key>
        <string>1</string>
    </dict>
</dict>
</plist>
